# Ruby/timekpr-web/timekpr-web.conf
description "Remote control of timekpr allowance via web"
author "Jeff Kowalski"
 
start on started mountall
stop on shutdown
 
script
  exec ruby /usr/sbin/timekpr-web
end script
